# debian/postinst
#!/bin/bash
# postinst script for jupiter
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
      if [ ! "$(grep jupiter /etc/group)" ]; then
        echo "Adding Jupiter group"
        /usr/sbin/groupadd jupiter
      fi
      echo -n "Adding users to Jupiter group: "
      for i in $(awk -F: '$6 ~ /\/home/ && $3 >= 1000 {print $1}' /etc/passwd);
      do
        echo -n $i" "
        usermod -G jupiter -a $i 2>&1 >/dev/null ||:
      done
      echo -e "\nAdding jupiter to sudoers"
      grep jupiter /etc/sudoers >/dev/null 2>&1 || echo "%jupiter ALL=NOPASSWD: /usr/lib/jupiter/scripts/bluetooth, /usr/lib/jupiter/scripts/camera, /usr/lib/jupiter/scripts/cpu-control, /usr/lib/jupiter/scripts/resolutions, /usr/lib/jupiter/scripts/rotate, /usr/lib/jupiter/scripts/touchpad, /usr/lib/jupiter/scripts/vga-out, /usr/lib/jupiter/scripts/wifi" >> /etc/sudoers
      echo "Setting permissions.."
      chown -R root:jupiter /usr/lib/jupiter
      chown -R root:root /etc/acpi/events

      if [ ! -d "/var/jupiter" ]; then
        mkdir /var/jupiter
      fi

      if [ -d "/var/jupiter" ]; then
	chown -R root:jupiter /var/jupiter
        chmod -R 775 /var/jupiter
      fi

      if [ -e "/usr/lib/jupiter/scripts/jupiter" ]; then
        /usr/lib/jupiter/scripts/jupiter ||:
      fi
      if [ -e "/usr/sbin/update-rc.d" ]; then
	if [ -e "/etc/init.d/ondemand" ]; then
	  echo "Removing startup links to ondemand CPU mode"
	  update-rc.d -f ondemand remove 2>&1 >/dev/null
	fi
      fi

      if [ -e "/etc/pm/sleep.d/00-jupiter-wifi" ]; then
        rm -f /etc/pm/sleep.d/00-jupiter-wifi 2>&1 >/dev/null ||:
      fi


      if [ -e "/etc/init.d/acpid" ]; then
        echo Restarting ACPI
        service acpid restart >/dev/null 2>&1 ||:
      else
        echo -e "\nPlease reboot your portable at your earliest convenience."
      fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
